# scripts/redis.conf
description "redis server"
author      "gadr90"

start on filesystem and started networking 
stop on shutdown

script
    exec /usr/local/bin/redis-server >> /var/log/redis.sys.log 2>&1
end script
